# docker/supervisor.conf
[supervisord]
nodaemon=true
pidfile=/home/application/supervisord.pid
logfile=/home/application/supervisord.log

[unix_http_server]
chown = application
file=/dev/shm/supervisor.sock
chmod=0700

[supervisorctl]
serverurl = unix:///dev/shm/supervisor.sock
username = root
password = docker

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface


[program:laravel-worker]
command=php /app/artisan queue:work --timeout 3600
user=application
numprocs=1
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
stdout_logfile=/app/storage/logs/laravel-worker.log
stdout_logfile_maxbytes=0

[include]
files = /opt/docker/etc/supervisor.d/*.conf
